# ALEX/scripts/systemd/alex-worker.service
[Unit]
Description=ALEX analyze worker
After=network.target

[Service]
Type=simple
User=alex
Group=alex
WorkingDirectory=/opt/alex/ALEX
Environment=PATH=/opt/alex/ALEX/.venv/bin:/usr/bin
EnvironmentFile=-/opt/alex/ALEX/.env
ExecStart=/opt/alex/ALEX/.venv/bin/python -m web.worker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
