# pm_test_spec_assistant/scripts/systemd/alex-worker.service
[Unit]
Description=ALEX analyze worker (production queue)
After=network.target alex-web.service

[Service]
Type=simple
User=alex
Group=alex
WorkingDirectory=/opt/alex/SyS_Spec_GenAI/pm_test_spec_assistant
Environment=PATH=/opt/alex/SyS_Spec_GenAI/pm_test_spec_assistant/.venv/bin:/usr/bin
ExecStart=/opt/alex/SyS_Spec_GenAI/pm_test_spec_assistant/.venv/bin/python -m web.worker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
